# debian/update-motd.d/10-help-text
#!/bin/sh

echo
echo "Welcome to Netrunner!"
echo " * Documentation:  http://www.netrunner-os.com/"

